# examples/shkd.service.template
[Unit]
Description=Simple HotKey Daemon

[Service]
Type=simple
ExecStart=BIN_PREFIX/shkd INPUT_DEVICE

[Install]
WantedBy=multi-user.target
