# examples/shkd.service.template
[Unit]
Description=Simple HotKey Daemon

[Service]
ExecStart=BIN_PREFIX/shkd %I

[Install]
WantedBy=multi-user.target
